# crontab
# Exec script hourly
0 * * * * /app/cron.py