# crontab
# Exec script hourly
0 * * * * /app/cron.py
59 * * * * find /app/ftpdata/* -type f -mtime +1 -delete;
