# crontab
# Exec script hourly
0 * * * * /app/cron.py
59 * * * * /app/scripts/cleaning.sh